# custom_script/package-list-backup.cron
0 2 * * * /home/webbacillus/.dotfiles/custom_script/package-list-backup.sh
